# deploy/douyin-downloader-api.service
[Unit]
Description=Douyin Downloader API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/douyin-downloader
Environment="DOUYIN_COOKIE=msToken=...; ttwid=...; odin_tt=...;"
ExecStart=/opt/douyin-downloader/.venv/bin/python -m uvicorn api_server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
